# run_sql.cron
30 23 * * 1-5 /usr/bin/python3 scraper_sql.py
